# modules/userland/files/smuxi-server.service
[Unit]
Description=Smuxi launcher

[Service]
ExecStart=/usr/bin/smuxi-server
User=smuxi

[Install]
WantedBy=multi-user.target
