# modules/userland/files/tunneler@.service
[Unit]
Description=Ssh tunnel remote script %I

[Service]
ExecStart=/usr/bin/tunnelscript %I

[Install]
WantedBy=multi-user.target
